# crontab
# 每小时执行一次脚本
*/10 * * * * cd /usr/src/app && /usr/local/bin/node --env-file=.env ./scripts/index.js >> /var/log/cron.log 2>&1
